# ansible/roles/ipv6_dns_sync/templates/au.hosteng.ipv6-dns-shutdown.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>au.hosteng.ipv6-dns-shutdown</string>

    <!-- Run the shutdown hook script which traps SIGTERM and runs DNS cleanup.
         The script stays alive via "tail -f /dev/null" and waits for SIGTERM,
         which launchd sends to all running daemons during system shutdown.
         This is more reliable than RunAtShutdown or SIGTERM handling in
         the Python watcher. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/ipv6_dns_shutdown_hook.sh</string>
    </array>

    <!-- Load at boot and keep running until shutdown -->
    <key>RunAtLoad</key>
    <true/>

    <!-- KeepAlive false — if the script exits unexpectedly don't restart it.
         We don't want it restarting mid-shutdown after cleanup has run. -->
    <key>KeepAlive</key>
    <false/>

    <!-- Give the script enough time to complete DNS cleanup before
         launchd forcibly kills it. -->
    <key>ExitTimeout</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/var/log/ipv6_dns_shutdown.out.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/ipv6_dns_shutdown.err.log</string>
</dict>
</plist>
