# etc/systemd/system/ipv6-dns-sync.service
[Unit]
Description=IPv6 DNS watcher and sync via nsupdate
Documentation=man:nsupdate(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ipv6_dns_watch_linux.sh
Restart=on-failure
RestartSec=5
User=root
Group=root
Nice=10
StandardOutput=append:/var/log/ipv6_sync.out
StandardError=append:/var/log/ipv6_sync.err

[Install]
WantedBy=multi-user.target
