# etc/systemd/system/ipv6-dns-sync.service
[Unit]
Description=Sync IPv6 AAAA/PTR records via nsupdate
Documentation=man:nsupdate(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ipv6_dns_sync.py --config-url https://a02.au/ipv6_sync_config/config.json
User=root
Group=root
Nice=10
StandardOutput=append:/var/log/ipv6_sync.out
StandardError=append:/var/log/ipv6_sync.err
# (IMPORTANT: removed the ExecStartPre line)
# No ExecStartPre here
