# etc/systemd/system/ipv6-dns-watch.service
[Unit]
Description=Watch IPv6 address changes and sync AAAA/PTR via nsupdate
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ipv6_dns_watch_linux.sh
User=root
Group=root
Restart=always
RestartSec=2
Nice=10
StandardOutput=append:/var/log/ipv6_watch.out
StandardError=append:/var/log/ipv6_watch.err

[Install]
WantedBy=multi-user.target
